# crontab
0 14 * * 5 /app/bin/weekly eval "Weekly.Sender.perform()"
