# crontab
0 21 * * 5 /app/bin/weekly eval "Weekly.Sender.perform()"
